# packages/scadere/scadere-0.2.4.tar.gz/scadere-0.2.4/contrib/scadere-check.timer
[Unit]
Description=Run scadere-check everyday

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=86400
FixedRandomDelay=true

[Install]
WantedBy=timers.timer
