# packages/scadere/scadere-0.2.4.tar.gz/scadere-0.2.4/contrib/scadere-listen.service
[Unit]
Description=Run scadere-listen
After=network.target

[Service]
ConditionPathExists=/var/lib/scadere/certificates
ExecStart=scadere-listen --title='TLS cert renewal' \
	/var/lib/scadere/certificates https://scadere.example :4433
Restart=on-failure

[Install]
WantedBy=default.target
