# packages/scadere/scadere-0.2.4.tar.gz/scadere-0.2.4/contrib/scadere-check.service
[Unit]
Description=Run scadere-check
After=network.target

[Service]
ExecStart=scadere-check --output=/var/lib/scadere/certificates \
	example.com example.net
