# packages/pimelon/pimelon-6.16.20.tar.gz/pimelon-6.16.20/pine/config/templates/supervisor.conf
; Notes:
; priority=1 --> Lower priorities indicate programs that start first and shut down last
; killasgroup=true --> send kill signal to child processes too

; graceful timeout should always be lower than stopwaitsecs to avoid orphan gunicorn workers.
[program:{{ pine_name }}-melon-web]
command={{ pine_dir }}/env/bin/gunicorn -b 127.0.0.1:{{ webserver_port }} -w {{ gunicorn_workers }} --max-requests {{ gunicorn_max_requests }} --max-requests-jitter {{ gunicorn_max_requests_jitter }} -t {{ http_timeout }} --graceful-timeout 30 melon.app:application --preload
priority=4
autostart=true
autorestart=true
stdout_logfile={{ pine_dir }}/logs/web.log
stderr_logfile={{ pine_dir }}/logs/web.error.log
stopwaitsecs=40
killasgroup=true
user={{ user }}
directory={{ sites_dir }}
startretries={{ supervisor_startretries }}

[program:{{ pine_name }}-melon-schedule]
command={{ pine_cmd }} schedule
priority=3
autostart=true
stdout_logfile={{ pine_dir }}/logs/schedule.log
stderr_logfile={{ pine_dir }}/logs/schedule.error.log
user={{ user }}
directory={{ pine_dir }}
startretries={{ supervisor_startretries }}

{% if not multi_queue_consumption %}
[program:{{ pine_name }}-melon-default-worker]
command={{ pine_cmd }} worker --queue default
priority=4
autostart=true
autorestart=true
stdout_logfile={{ pine_dir }}/logs/worker.log
stderr_logfile={{ pine_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=1560
directory={{ pine_dir }}
killasgroup=true
numprocs={{ background_workers }}
process_name=%(program_name)s-%(process_num)d
startretries={{ supervisor_startretries }}
{% endif %}

[program:{{ pine_name }}-melon-short-worker]
command={{ pine_cmd }} worker --queue short{{',default' if multi_queue_consumption else ''}}
priority=4
autostart=true
autorestart=true
stdout_logfile={{ pine_dir }}/logs/worker.log
stderr_logfile={{ pine_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=360
directory={{ pine_dir }}
killasgroup=true
numprocs={{ background_workers }}
process_name=%(program_name)s-%(process_num)d
startretries={{ supervisor_startretries }}

[program:{{ pine_name }}-melon-long-worker]
command={{ pine_cmd }} worker --queue long{{',default,short' if multi_queue_consumption else ''}}
priority=4
autostart=true
autorestart=true
stdout_logfile={{ pine_dir }}/logs/worker.log
stderr_logfile={{ pine_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=1560
directory={{ pine_dir }}
killasgroup=true
numprocs={{ background_workers }}
process_name=%(program_name)s-%(process_num)d
startretries={{ supervisor_startretries }}

{% for worker_name, worker_details in workers.items() %}
[program:{{ pine_name }}-melon-{{ worker_name }}-worker]
command={{ pine_cmd }} worker --queue {{ worker_name }}
priority=4
autostart=true
autorestart=true
stdout_logfile={{ pine_dir }}/logs/worker.log
stderr_logfile={{ pine_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs={{ worker_details["timeout"] }}
directory={{ pine_dir }}
killasgroup=true
numprocs={{ worker_details["background_workers"] or background_workers }}
process_name=%(program_name)s-%(process_num)d
startretries={{ supervisor_startretries }}
{% endfor %}


{% if not skip_redis %}
[program:{{ pine_name }}-redis-cache]
command={{ redis_server }} {{ redis_cache_config }}
priority=1
autostart=true
autorestart=true
stdout_logfile={{ pine_dir }}/logs/redis-cache.log
stderr_logfile={{ pine_dir }}/logs/redis-cache.error.log
user={{ user }}
directory={{ sites_dir }}
startretries={{ supervisor_startretries }}

[program:{{ pine_name }}-redis-queue]
command={{ redis_server }} {{ redis_queue_config }}
priority=1
autostart=true
autorestart=true
stdout_logfile={{ pine_dir }}/logs/redis-queue.log
stderr_logfile={{ pine_dir }}/logs/redis-queue.error.log
user={{ user }}
directory={{ sites_dir }}
startretries={{ supervisor_startretries }}
{% endif %}

{% if node %}
[program:{{ pine_name }}-node-socketio]
command={{ node }} {{ pine_dir }}/apps/melon/socketio.js
priority=4
autostart=true
autorestart=true
stdout_logfile={{ pine_dir }}/logs/node-socketio.log
stderr_logfile={{ pine_dir }}/logs/node-socketio.error.log
user={{ user }}
directory={{ pine_dir }}
startretries={{ supervisor_startretries }}
{% endif %}

[group:{{ pine_name }}-web]
programs={{ pine_name }}-melon-web {%- if node -%} ,{{ pine_name }}-node-socketio {%- endif%}


{% if multi_queue_consumption %}

[group:{{ pine_name }}-workers]
programs={{ pine_name }}-melon-schedule,{{ pine_name }}-melon-short-worker,{{ pine_name }}-melon-long-worker{%- for worker_name in workers -%},{{ pine_name }}-melon-{{ worker_name }}-worker{%- endfor %}

{% else %}

[group:{{ pine_name }}-workers]
programs={{ pine_name }}-melon-schedule,{{ pine_name }}-melon-default-worker,{{ pine_name }}-melon-short-worker,{{ pine_name }}-melon-long-worker{%- for worker_name in workers -%},{{ pine_name }}-melon-{{ worker_name }}-worker{%- endfor %}

{% endif %}

{% if not skip_redis %}
[group:{{ pine_name }}-redis]
programs={{ pine_name }}-redis-cache,{{ pine_name }}-redis-queue
{% endif %}
