# packages/pimelon/pimelon-6.16.20-py3-none-any.whl/pine/config/templates/systemd/pimelon-node-socketio.service
[Unit]
After={{ pine_name }}-melon-web.service
Description="{{ pine_name }}-node-socketio"
PartOf={{ pine_name }}-web.target

[Service]
User={{ user }}
Group={{ user }}
Restart=always
ExecStart={{ node }} {{ pine_dir }}/apps/melon/socketio.js
StandardOutput=file:{{ pine_dir }}/logs/node-socketio.log
StandardError=file:{{ pine_dir }}/logs/node-socketio.error.log
WorkingDirectory={{ pine_dir }}
