# packages/pimelon/pimelon-6.16.20.tar.gz/pimelon-6.16.20/pine/config/templates/systemd/pimelon.target
[Unit]
After=network.target
Requires={{ pine_name }}-web.target {{ pine_name }}-workers.target {{ pine_name }}-redis.target

[Install]
WantedBy=multi-user.target
