# packages/pimelon/pimelon-6.16.20.tar.gz/pimelon-6.16.20/pine/config/templates/systemd/pimelon-web.target
[Unit]
After=network.target
Wants={{ pine_name }}-melon-web.service {{ pine_name }}-node-socketio.service

[Install]
WantedBy=multi-user.target
