# packages/pimelon/pimelon-6.16.20.tar.gz/pimelon-6.16.20/pine/config/templates/systemd/pimelon-melon-short-worker.service
[Unit]
Description="{{ pine_name }}-melon-short-worker %I"
PartOf={{ pine_name }}-workers.target

[Service]
User={{ user }}
Group={{ user }}
Restart=always
ExecStart={{ pine_cmd }} worker --queue short
StandardOutput=file:{{ pine_dir }}/logs/worker.log
StandardError=file:{{ pine_dir }}/logs/worker.error.log
WorkingDirectory={{ pine_dir }}
